# cPodRouter/systemd/dyndns.timer
[Unit]
Description=Update public DNS entry for this network

[Timer]
OnBootSec=1min
OnCalendar=*:0/1

[Install]
WantedBy=basic.target
